# devpi/devpi.service.in
[Unit]
Requires=network-online.target
After=network-online.target

Type=forking
PIDFile=${DEVPI_HOME}/.xproc/devpi-server/xprocess.PID
Restart=always

ExecStart=/usr/local/bin/devpi-server --port 4040 --serverdir ${DEVPI_HOME} --start
ExecStop=/usr/local/bin/devpi-server --port 4040 --serverdir ${DEVPI_HOME} --stop

User=${DEVPI_USER}
